# distros/pmos-fbkeyboard/overrides/etc/init.d/eg25
#!/sbin/openrc-run

description="Starting EG25 WWAN module"

depend()
{
	before ofono
}

start()
{
	ebegin "Enabling EG25 WWAN module"
	echo 1 > /sys/class/modem-power/modem-power/device/powered
}

stop()
{
	ebegin "Disabling EG25 WWAN module"
	echo 0 > /sys/class/modem-power/modem-power/device/powered
}
